# mongrel2-main/supervisord.conf
[supervisord]
nodaemon=true

[program:procer]
command=/bin/bash -c "cd /deployment && procer /deployment/profiles /deployment/procer.pid"
